# systemd/telegram-ai-cli-bridge.service
[Unit]
Description=Telegram AI CLI Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/longtail1019/telecodex/bridge.py
WorkingDirectory=/home/longtail1019/telecodex
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
